# isrsrv-timer-2.service
[Unit]
Description=IsrSrv Script Timer 2 Service

[Service]
Type=oneshot
ExecStart=/usr/bin/isrsrv-script timer_two
