# isrsrv-timer-1.service
[Unit]
Description=IsrSrv Script Timer 1 Service

[Service]
Type=oneshot
ExecStart=/usr/bin/isrsrv-script timer_one
